# packages/MyNewspaper/mynewspaper-4.0.tar.gz/mynewspaper-4.0/tools/mynewspaper.cron
# Method 1: by url: HTTP
11,41 * * * * /usr/bin/wget -q -O /dev/null http://127.0.0.1:8087/update

# Method 2: by url: HTTPS and protected by user&password
#11,41 * * * * /usr/bin/wget -q -O /dev/null --no-check-certificate --http-user=username --http-password=password https://mynewspaper.url/update
