# packages/pywizard/pywizard-0.5.12.tar.gz/pywizard-0.5.12/pywizard/templates/upstart-agent.conf
start on startup
stop on shutdown

script
    exec pywizard-agent %ARGS% > /var/log/pywizard-agent.log 2>&1
end script

respawn