# packages/pywizard/pywizard-0.5.12.tar.gz/pywizard-0.5.12/pywizard/templates/upstart-server.conf
start on startup
stop on shutdown

script
    exec pywizard-server > /var/log/pywizard-server.log 2>&1
end script