# packages/pscripts/pscripts-0.1.160.tar.gz/pscripts-0.1.160/config/update_external_ip.service
[Unit]
Description=Update your dynamic IP address with a Dynamic DNS service

[Service]
Type=forking
ExecStart=/usr/bin/update_external_ip --start
ExecStop=/usr/bin/update_external_ip --stop
PIDFile=/var/run/update_external_ip.pid
Restart=on-abort

[Install]
WantedBy=multi-user.target